# crontab
*/2 * * * *   root    /etc/cron.custom/minecraft-keepalive
